# execution/daemons/sylvia/com.ateles.sylvia.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.sylvia</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/repos/ateles/.venv/bin/python3</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/sylvia/sylvia.py</string>
    </array>

    <!-- Run daily at 07:15 UTC (09:15 Madrid summer / 08:15 winter) — after Cotinga briefing -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/sylvia.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/sylvia.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/markmhendrickson</string>
        <key>PATH</key>
        <string>/Users/markmhendrickson/.local/bin:/Users/markmhendrickson/.nvm/versions/node/v20.20.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NEOTOMA_BASE_URL</key>
        <string>http://localhost:3180</string>
    </dict>
</dict>
</plist>
