# execution/scripts/com.markmhendrickson.formica.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.markmhendrickson.formica</string>

    <key>ProgramArguments</key>
    <array>
        <string>@ATELES_ROOT@/execution/scripts/run_formica_launchd.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>@ATELES_ROOT@</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>@ATELES_ROOT@/data/logs/formica_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>@ATELES_ROOT@/data/logs/formica_launchd.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
