# execution/scripts/com.finances.tunnel-url-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.finances.tunnel-url-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/Projects/personal/execution/venv/bin/python</string>
        <string>/Users/markmhendrickson/Projects/personal/execution/scripts/monitor_tunnel_url.py</string>
        <string>--interval</string>
        <string>30</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/markmhendrickson/Projects/personal</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/markmhendrickson/Projects/personal/data/logs/tunnel_url_monitor.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/markmhendrickson/Projects/personal/data/logs/tunnel_url_monitor.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>






