# execution/scripts/com.cloudflare.twilio-sms-webhook-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cloudflare.twilio-sms-webhook-tunnel</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>twilio-sms-webhook</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/markmhendrickson/Projects/personal</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/markmhendrickson/Projects/personal/data/logs/cloudflare_twilio_sms_tunnel.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/markmhendrickson/Projects/personal/data/logs/cloudflare_twilio_sms_tunnel.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>

