# execution/daemons/neotoma-agent/com.ateles.neotoma-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.neotoma-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/repos/ateles/.venv/bin/python3</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/neotoma-agent/neotoma_agent.py</string>
    </array>

    <!-- Load at login -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/neotoma-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/neotoma-agent.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/markmhendrickson</string>
        <key>PATH</key>
        <string>/Users/markmhendrickson/.local/bin:/Users/markmhendrickson/.nvm/versions/node/v20.20.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NEOTOMA_BASE_URL</key>
        <string>http://localhost:3180</string>
    </dict>
</dict>
</plist>
