# execution/daemons/anthus/com.ateles.anthus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.anthus</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/repos/ateles/.venv/bin/python3</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/anthus/anthus.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/markmhendrickson/repos/ateles/execution/daemons/anthus</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/markmhendrickson/.nvm/versions/node/v22.15.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NEOTOMA_BASE_URL</key>
        <string>https://neotoma.markmhendrickson.com</string>
        <key>NEOTOMA_BEARER_TOKEN</key>
        <string>__NEOTOMA_BEARER_TOKEN__</string>
        <key>NEOTOMA_SSE_SUBSCRIPTION_ID_ANTHUS</key>
        <string>__NEOTOMA_SSE_SUBSCRIPTION_ID_ANTHUS__</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.ateles.anthus.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.ateles.anthus.err.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
