# execution/daemons/gorilla/com.ateles.gorilla.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.gorilla</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/gorilla/gorilla.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/markmhendrickson/repos/ateles/execution/daemons/gorilla</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.ateles.gorilla.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.ateles.gorilla.err.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
