# execution/daemons/phoenicurus-release/com.ateles.phoenicurus-prepare.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.phoenicurus-prepare</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/repos/ateles/.venv/bin/python3</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/phoenicurus-release/prepare.py</string>
    </array>

    <!-- Mon-Thu at 07:00 local. Weekday: 0=Sun .. 6=Sat. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>1</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>2</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>3</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key><integer>4</integer>
            <key>Hour</key><integer>7</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/markmhendrickson</string>
        <key>PATH</key>
        <string>/Users/markmhendrickson/.nvm/versions/node/v20.20.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NEOTOMA_BASE_URL</key>
        <string>http://localhost:3180</string>
        <key>NEOTOMA_REPO_ROOT</key>
        <string>/Users/markmhendrickson/repos/neotoma</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/phoenicurus-release.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/phoenicurus-release.log</string>
</dict>
</plist>
