# execution/daemons/tyto/com.ateles.tyto.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.tyto</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/repos/personal/execution/venv/bin/python</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/tyto/tyto.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/markmhendrickson/.local/bin:/Users/markmhendrickson/.nvm/versions/node/v20.20.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/markmhendrickson</string>
        <key>VIRTUAL_ENV</key>
        <string>/Users/markmhendrickson/repos/personal/execution/venv</string>
        <key>PYTHONPATH</key>
        <string>/Users/markmhendrickson/repos/ateles:/Users/markmhendrickson/repos/personal/execution/venv/lib/python3.11/site-packages</string>
        <key>NEOTOMA_BASE_URL</key>
        <string>https://neotoma.markmhendrickson.com</string>
        <key>TYTO_SCREENSHOTS_DIR</key>
        <string>/Users/markmhendrickson/Desktop/Screenshots</string>
        <key>TYTO_POLL_INTERVAL</key>
        <string>10</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/tyto.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/markmhendrickson/Library/Logs/ateles/tyto.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
