# execution/daemons/castor/com.ateles.castor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.castor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/repos/ateles/execution/venv/bin/python</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/castor/castor.py</string>
    </array>

    <!-- Load at login -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Restart if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.ateles.castor.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.ateles.castor.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/markmhendrickson</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
