# execution/daemons/apus/com.ateles.apus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.apus</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/apus/apus.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>NEOTOMA_BASE_URL</key>
        <string>https://neotoma.markmhendrickson.com</string>
        <key>ATELES_REPO_PATH</key>
        <string>/Users/markmhendrickson/repos/ateles</string>
        <key>ATELES_PRIVATE_REPO_PATH</key>
        <string>/Users/markmhendrickson/repos/ateles-private</string>
        <key>ATELES_AGENT_GIT_NAME</key>
        <string>ateles-agent</string>
        <key>ATELES_AGENT_GIT_EMAIL</key>
        <string>ateles-agent@users.noreply.github.com</string>
        <key>APUS_WEBHOOK_SECRET</key>
        <string>apus-neotoma-webhook-2026</string>
        <key>APUS_PORT</key>
        <string>8741</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.ateles.apus.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.ateles.apus.err.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
