# execution/daemons/morning-brief/com.ateles.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ateles.morning-brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/markmhendrickson/repos/ateles/.venv/bin/python3</string>
        <string>/Users/markmhendrickson/repos/ateles/execution/daemons/morning-brief/morning-brief.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/com.ateles.morning-brief.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.ateles.morning-brief.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/markmhendrickson</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin:/Users/markmhendrickson/.nvm/versions/node/v22.15.0/bin</string>
    </dict>
</dict>
</plist>
